# OpenClinica/vendor/ndar/OpenClinica-2.0/conf/update_data_warehouse.cron
#!/bin/sh



/usr/local/pgsql/bin/psql -U clinica -d openclinica -f /var/OpenClinica-1.0/conf/pg_data_warehouse.sql
